# env/cron.d/scrap
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
# m h dom mon dow user<>command
0,15,30,45 * * * * root python /var/www/scrap/scrap.py "https://www.olx.pl/motoryzacja/samochody/volkswagen/caravelle/?search%5Bfilter_float_year%3Ato%5D=1990"
0,15,30,45 * * * * root python /var/www/scrap/scrap.py "https://www.olx.pl/motoryzacja/samochody/volkswagen/multivan/?search%5Bfilter_float_year%3Ato%5D=1990"
0,15,30,45 * * * * root python /var/www/scrap/scrap.py "https://www.olx.pl/motoryzacja/samochody/volkswagen/transporter/?search%5Bfilter_float_year%3Ato%5D=1990"
0,15,30,45 * * * * root python /var/www/scrap/scrap.py "https://www.olx.pl/motoryzacja/dostawcze-ciezarowe/q-t2/?search%5Bfilter_float_price%3Ato%5D=15000&search%5Bfilter_float_year%3Ato%5D=1991"
0,15,30,45 * * * * root python /var/www/scrap/scrap.py "https://www.olx.pl/motoryzacja/dostawcze-ciezarowe/q-t3/?search%5Bfilter_float_price%3Ato%5D=15000&search%5Bfilter_float_year%3Ato%5D=1991"
0,15,30,45 * * * * root python /var/www/scrap/scrap.py "https://www.olx.pl/motoryzacja/samochody/volkswagen/garbus/pomorskie/?search%5Bfilter_float_price%3Ato%5D=15000&search%5Bfilter_float_year%3Ato%5D=1990"
0,15,30,45 * * * * root python /var/www/scrap/scrap.py "https://www.olx.pl/motoryzacja/samochody/pozostale/pomorskie/?search%5Bfilter_float_price%3Ato%5D=5000&search%5Bfilter_float_year%3Ato%5D=1977"





